# fan-control.service
[Unit]
Description=Fan Control Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/fan-control
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target 